# scripts/systemd/cc-frontend.service
[Unit]
Description=CommandCenter Frontend (Next.js Production)
After=network.target cc-backend.service

[Service]
Type=simple
WorkingDirectory=/home/rohith/desktop/CommandCenter/frontend
Environment="PATH=/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="HOSTNAME=0.0.0.0"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=append:/home/rohith/desktop/CommandCenter/logs/frontend.log
StandardError=append:/home/rohith/desktop/CommandCenter/logs/frontend.log

[Install]
WantedBy=default.target
