# claude-rl-agent/engine/claude-rl-capture.service
[Unit]
Description=Claude Code RL Capture Daemon
After=network.target

[Service]
Type=simple
User=rohith
WorkingDirectory=/home/rohith/desktop/CommandCenter/claude-rl-agent
ExecStart=/home/rohith/.pyenv/versions/3.11.9/bin/python3 /home/rohith/desktop/CommandCenter/claude-rl-agent/engine/auto_capture_daemon.py start
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/rohith/desktop/CommandCenter/claude-rl-agent/logs/daemon.log
StandardError=append:/home/rohith/desktop/CommandCenter/claude-rl-agent/logs/daemon-error.log

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"

[Install]
WantedBy=multi-user.target
